# services/ssh/run
#!/bin/sh

exec /usr/sbin/sshd -D -o StrictModes=no 2> /var/log/sshd.log